# systemd/btrfs-maint@.timer
[Unit]
Description=btrfs scrub timer for %f

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
