# systemd/btrfs-maint@.service
[Unit]
Description=btrfs scrub service for %f

[Service]
Type=simple
ExecStart=/var/JohnstonNetwork/tasks/btrfs-maint %f
