# systemd/efi-mirror.timer
[Unit]
Description=EFI partition mirroring timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
