# systemd/lvm-maint@.service
[Unit]
Description=lvm RAID maintenance service for %I

[Service]
Type=simple
ExecStart=/var/JohnstonNetwork/tasks/lvm-maint %I
