# systemd/lynis.timer
[Unit]
Description=Periodically run Lynis security auditing

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
