# systemd/update-ad-dns.timer
[Unit]
Description=Update DNS registration via SAMBA

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
