# systemd/update-ad-dns.service
[Unit]
Description=Update DNS registration via SAMBA

[Service]
Type=simple
ExecStart=/var/JohnstonNetwork/tasks/update-ad-dns
