# systemd/zpool-maint@.timer
[Unit]
Description=ZFS zpool scrub timer for %I

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
