# systemd/zpool-maint@.service
[Unit]
Description=ZFS zpool scrub service for %I

[Service]
Type=simple
ExecStart=/sbin/zpool scrub %I
